# etc/xdg/autostart/emmi-update.desktop
[Desktop Entry]
Name=Emmi Update
Comment=Emmi Update Manager
Icon=stock_lock
Exec=emmi-update
Terminal=false
Type=Application
Categories=Utility;
